# node_exp/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=root
ExecStart={{ node_exporter_bin_path }}
Restart=always

[Install]
WantedBy=multi-user.target
